# deploy/systemd/safe-kiosk-people-agent.slice
[Unit]
Description=Safe Kiosk People Agent resource slice
[Slice]
MemoryHigh=1536M
MemoryMax=2048M
